# scripts/supervisord-prax.conf
; Prax container supervisor — keeps Qdrant, Neo4j, TeamWork, optional
; ngrok, and the Flask watchdog all running.  Replaces the prior
; entrypoint that forked everything from a single bash script (children
; got reparented to whichever process was PID 1, with no restart logic).
;
; Process tree at runtime:
;   PID 1: supervisord
;     ├─ qdrant     (priority 10) — vector store on :6333
;     ├─ neo4j      (priority 15) — graph store on :7474 / :7687
;     ├─ teamwork   (priority 20) — FastAPI web UI on :8000
;     ├─ ngrok      (priority 25) — only if NGROK_AUTHTOKEN is set; otherwise sleeps
;     └─ watchdog   (priority 30) — supervises the Flask app on :5001
;
; Every program logs to /dev/fd/1 + /dev/fd/2 so `docker compose logs prax`
; keeps showing what it did before.  All have autorestart=true, so a
; crashed Neo4j / TeamWork / Flask comes back without a container restart.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ── Vector store ────────────────────────────────────────────────────

[program:qdrant]
command=/usr/local/bin/qdrant
priority=10
autostart=true
autorestart=true
startretries=10
startsecs=3
environment=QDRANT__STORAGE__STORAGE_PATH="/data/qdrant",QDRANT__SERVICE__HTTP_PORT="6333",QDRANT__SERVICE__GRPC_PORT="6334"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Graph store ─────────────────────────────────────────────────────

[program:neo4j]
command=/opt/neo4j/bin/neo4j console
priority=15
autostart=true
autorestart=true
startretries=10
startsecs=10
; Heap caps keep Neo4j under ~450 MB so the bundled stack fits on small
; hosts.  Auth + plugin list come in via the entrypoint's exported env.
environment=NEO4J_HOME="/opt/neo4j",NEO4J_server_memory_heap_initial__size="128m",NEO4J_server_memory_heap_max__size="256m",NEO4J_server_memory_pagecache_size="64m",NEO4J_dbms_usage__report__enabled="false"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── TeamWork (FastAPI + React UI) ───────────────────────────────────

[program:teamwork]
command=/app/scripts/teamwork-launch.sh
priority=20
autostart=true
autorestart=true
startretries=10
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Optional ngrok (Twilio + opt-in shares) ─────────────────────────
; Wrapper checks NGROK_AUTHTOKEN; if unset it sleeps forever so
; supervisord doesn't busy-loop trying to restart a no-op program.

[program:ngrok]
command=/app/scripts/ngrok-launch.sh
priority=25
autostart=true
autorestart=true
startretries=5
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

; ── Flask app + its own watchdog ─────────────────────────────────────
; The watchdog still runs Flask under itself for fast restarts on
; in-process exceptions; supervisord wraps the watchdog so a crash of
; the supervisor itself is recovered too.

[program:watchdog]
command=/app/scripts/watchdog-launch.sh
priority=30
autostart=true
autorestart=true
startretries=10
startsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
